# etc/gnome/autostart/emacs.desktop
[Desktop Entry]
Type=Application
Name=Emacs
Exec=emacs
