# etc/gnome/autostart/init.desktop
[Desktop Entry]
Type=Application
Name=Init
Exec=/home/taylan/bin/xinitrc
