# etc/systemd/bitlbee.service
[Unit]
Description=BitlBee Server

[Service]
ExecStart=/bin/sh -c 'exec bitlbee -Dn -c "$HOME_CONFIG"/bitlbee/bitlbee.conf'
