# etc/systemd/emacs.service
[Unit]
Description=Emacs Daemon

[Service]
Type=forking
ExecStart=/bin/sh -c 'emacs-daemon'
