# etc/systemd/freedns.service
[Unit]
Description=FreeDNS Updater

[Service]
Restart=always
RestartSec=5m
ExecStart=/bin/sh -c 'exec freedns-updater'
