# etc/systemd/freenet.service
[Unit]
Description=Freenet Node

[Service]
Type=forking
ExecStart=/bin/sh -c '~/.junk/Freenet/run.sh start'
