# etc/systemd/isync.timer
[Unit]
Description=IMAP Synchronizer

[Timer]
Unit=isync.service
OnCalendar=*-*-* *:00/10:00
