# etc/systemd/polipo.service
[Unit]
Description=Polipo HTTP Proxy

[Service]
ExecStart=/bin/sh -c 'exec polipo diskCacheRoot= logFile=/dev/null'
