# etc/systemd/rtorrent.service
[Unit]
Description=rTorrent

[Service]
Type=forking
ExecStart=/bin/sh -c 'tmux new -d -s rtorrent rtorrent'
