# etc/systemd/sshfs-t61.service
[Unit]
Description=Mount / from T61

[Service]
ExecStart=/bin/sh -c 'exec sshfs -f root@t61:/ ~/t61'
