# etc/systemd/sshfs-media.service
[Unit]
Description=Mount ~/media from T61

[Service]
ExecStart=/bin/sh -c 'exec sshfs -f t61:media ~/media'
