# etc/systemd/tor.service
[Unit]
Description=Tor Node

[Service]
ExecStart=/bin/sh -c 'exec tor -f "$HOME_CONFIG"/tor/rc RunAsDaemon 0'
