# etc/systemd/xorg.service
[Unit]
Description=X.Org Server

[Service]
ExecStart=/bin/sh -c 'exec xinit -- vt1'
